# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Start cron
# ==============================================================================

bashio::log.info "Starting cron..."

exec crond -f -l 2
