# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Start nginx
# ==============================================================================

bashio::log.info "Starting nginx..."

nginx -g "daemon off;"
