# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/init-wg++/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Updates WG++ channels in WG++ ini
# ==============================================================================
declare wg_channel
declare country
declare site
declare channel
declare filename
declare wg_config

bashio::log.info "Init WG++..."
# create WG++ guide folder in addon config folder
if ! bashio::fs.directory_exists "/config/tvheadend/wg++"; then
  bashio::log.info "Creating WG++ folder in config folder..."
  mkdir /config/tvheadend/wg++
fi

expectedParentDir="/wg/.wg++/siteini.pack"

function verifySiteFile() {
  local fullpath

  fullpath="$expectedParentDir"/"$1"/"$2".channels.xml
  if ! bashio::fs.file_exists "$fullpath"; then
    bashio::log.fatal
    bashio::log.fatal "WG Channel Configuration is incorrect for country $1 site $2"
    bashio::log.fatal "File $1/$2.channels.xml does not exist. Only use contries and sites"
    bashio::log.fatal "(case-sensitive) from the WG++ siteini.pack folder"
    bashio::exit.nok
  fi

  rpath=$(realpath "$fullpath")
  if [[ $rpath != ${expectedParentDir}* ]]; then
    bashio::log.fatal
    bashio::log.fatal "WG Channel Configuration is incorrect for country $1 site $2"
    bashio::log.fatal "File $1/$2.channels.xml is traversing paths. Only use contries and sites"
    bashio::log.fatal "(case-sensitive) from the WG++ siteini.pack folder"
    bashio::log.fatal
    bashio::exit.nok
  fi
  echo "$fullpath"
}

if bashio::config.has_value 'wg_config' && bashio::config.has_value 'wg_channels'; then
  bashio::log.fatal
  bashio::log.fatal "Both WG++ channels and a WG++ config xml are set. Only one can be set"
  bashio::log.fatal
  bashio::exit.nok
fi

# remove dummy channel from initial config xml
xmlstarlet ed -L -d "/settings/channel" /wg/.wg++/WebGrab++.config.xml

# remove license from initial config xml causing "Index was outside the bounds of the array." error
xmlstarlet ed -L -d "/settings/license" /wg/.wg++/WebGrab++.config.xml

if bashio::config.has_value 'wg_config'; then
  bashio::log.info "Using supplied WG++ config..."
  # Overwrite config file from addon config
  wg_config=$(bashio::config 'wg_config')
  echo "${wg_config}" > /wg/.wg++/WebGrab++.config.xml
  xmlstarlet val -q /wg/.wg++/WebGrab++.config.xml || bashio::exit.nok "Invalid WG++ config supplied"
  xmlstarlet fo /wg/.wg++/WebGrab++.config.xml > /wg/.wg++/tmp.xml && mv /wg/.wg++/tmp.xml /wg/.wg++/WebGrab++.config.xml
fi

# Add Channels into WG++ ini
if bashio::config.has_value 'wg_channels'; then
  bashio::log.info "Updating WG++ Channels..."
  for wg_channel in $(bashio::config 'wg_channels'); do
    country=$(bashio::jq "$wg_channel" '.country')
    site=$(bashio::jq "$wg_channel" '.site')
    channel=$(bashio::jq "$wg_channel" '.channel')
    bashio::log.info "Adding channel country: ${country} site: ${site} channel: ${channel} to wg++ config"

    filename=$(verifySiteFile "$country" "$site")
    if [[ ! $(xmlstarlet sel -t -v "//channel[@site_id='$channel']" "$filename") ]]; then
      bashio::log.fatal "Channel $channel does not exist in channel file of $site"
      bashio::exit.nok
    fi

    # run XSLT to copy over channel
    xmlstarlet tr /wg/tr.xsl -s filename="$filename" -s channel="$channel" /wg/.wg++/WebGrab++.config.xml | xmlstarlet fo >/wg/.wg++/tmp.xml &&
      mv /wg/.wg++/tmp.xml /wg/.wg++/WebGrab++.config.xml

  done
fi
